# pi-config/printer.cron
* * * * * cd /home/tim/tims-printer && scripts/printer.sh
